# vinotes/config/upstart-vinotes-api.conf
description "Gunicorn server for Vinotes API"

start on net-device-up
stop on shutdown
respawn

setuid rc
exec /home/rc/projects/vinotes-api/vinotes/bin/start_server.sh
